# deploy/worldcup.service
[Unit]
Description=2026 DG World Cup V3 API
After=network.target

[Service]
Type=simple
User=worldcup
Group=worldcup
WorkingDirectory=/opt/worldcup
EnvironmentFile=/etc/worldcup.env
ExecStart=/usr/bin/node server/app.js
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/worldcup

[Install]
WantedBy=multi-user.target
